# deploy/solana-blacklist-api.service
[Unit]
Description=Solana Blacklist API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/solana-blacklist
ExecStart=/var/www/solana-blacklist/target/release/api
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
